# systemd/server-cleanup.service
[Unit]
Description=Weekly server cleanup (Docker images, apt cache, journal logs)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/server-cleanup.sh
Nice=19
IOSchedulingClass=idle
